# supervisor-app.conf
[program:app]
command = /usr/local/bin/uwsgi --ini /rest-framework-tutorial/uwsgi.ini
